# cookbooks/node_app/templates/node.service.erb
[Unit]
Description=node-<%= @name %>
Requires=<%= @requires.join(' ') %>
After=<%= @requires.join(' ') %>

[Service]
Type=simple
User=<%= @user %>
Group=<%= @group %>
WorkingDirectory=<%= @base_dir %>/current
SyslogIdentifier=node-<%= @name %>
StandardOutput=syslog
StandardError=syslog
<%- @env.each do |env_variable| %>
Environment=<%= env_variable['name'] %>=<%= env_variable['value'] %>
<%- end %>
<%- @pre.each do |pre| %>
ExecStartPre=<%= pre %>
<%- end %>

ExecStart=<%= @node_path %> <%= @entrypoint %>

ExecStop=/bin/kill -s QUIT $MAINPID

RestartSec=1
Restart=on-failure

[Install]
WantedBy=multi-user.target
